# _generated/smock.service
# @axiom: infrastructure.md#systemd-service
[Unit]
Description=smock (well app)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/smock
ExecStart=/opt/smock/bin/smock
Restart=on-failure
RestartSec=3

EnvironmentFile=/opt/smock/.env

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/smock
PrivateTmp=yes


[Install]
WantedBy=multi-user.target
# /@axiom: infrastructure.md#systemd-service
